# virtualenv/myapp.conf
start on startup
expect daemon
#source /root/.env/bin/activate
script
su - root -c "source /root/.env/bin/activate"
su - root -c "/usr/bin/python2.7 /root/jinja/j.py"
end script

